# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/hotplug.d/usb/00_wwan.sh
#!/bin/ash

[ "$ACTION" = add ] || [ "$ACTION" = remove ] && [ "$DEVTYPE" = usb_device ] || exit 0

dname="${DEVPATH##*/}"

ubus send atsd.usb "{'action':'$ACTION','devpath':'$dname','from':'usb'}"

exit 0
