# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/hotplug.d/iface/99-adjust-usb-iface-rps
[ ifup = "$ACTION" ] && {
	[ -n "$DEVICE" ] && {
		if [ "$INTERFACE" = "cpe" ]; then
			echo 3 > /sys/class/net/${DEVICE}/queues/rx-0/rps_cpus
		fi
	}
}
