# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/access_ctl.list
/etc/init.d/access_ctl
/etc/config/access_ctl
/usr/bin/access_ctl.sh
/etc/hotplug.d/iface/99-access_ctl
