# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/appcenter.list
/etc/config/appcenter
/etc/init.d/appcenter
/usr/sbin/appcenter
