# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/atsd.list
/etc/init.d/atsd
/usr/sbin/atsd_cli
/usr/sbin/atsd
