# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/atserver-sniffer.list
/etc/init.d/atserver-sniffer
/usr/sbin/atserver-sniffer
