# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/cellular_record.list
/etc/config/cellular_record
/etc/init.d/cellular_record
/usr/sbin/cellular_record
