# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/cellular_init.list
/etc/config/cellular_init
/usr/sbin/cellular_init
/etc/init.d/cellular_init
