# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/cimd.list
/usr/sbin/cimd
/etc/init.d/cimd
/etc/config/cimd
