# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/cloudd.list
/usr/lib/cloudd/wifiset
/usr/lib/cloudd/init_only
/usr/lib/cloudd/cpestatus
/usr/bin/cloudd_reload_service
/usr/lib/cloudd/firmware
/etc/config/cloudd_acl
/usr/lib/cloudd/wifistatus
/usr/lib/cloudd/synctime
/usr/lib/cloudd/linkstatus
/usr/lib/cloudd/config
/usr/lib/cloudd/limitstatus
/usr/lib/lua/cloudd/file.lua
/usr/lib/cloudd/stastatus
/usr/lib/cloudd/gencfg
/usr/lib/cloudd/connection
/usr/lib/cloudd/traffic
/usr/lib/cloudd/report
/usr/lib/cloudd/file
/usr/lib/cloudd/develop
/usr/lib/cloudd/batstatus
/usr/lib/cloudd/status
/usr/lib/lua/cloudd/debug.lua
/usr/lib/cloudd/appstore
/usr/lib/cloudd/terminal
/usr/lib/cloudd/init-tdredcap
/usr/bin/cloudd_sysupgrade
/usr/lib/cloudd/reset
/usr/lib/cloudd/group_id
/usr/lib/lua/cloudd/api.lua
/usr/lib/cloudd/wanstatus
/usr/lib/cloudd/mosq_reload
/usr/bin/cloudd_cmd
/etc/config/cloudd_cli
/usr/lib/cloudd/upgrade
/usr/lib/cloudd/bdinfo
/usr/lib/cloudd/password
/usr/lib/cloudd/simswitch
/usr/lib/cloudd/command
/etc/init.d/cloudd
/usr/bin/cloudd
/usr/lib/lua/cloudd/sync_cloudd.lua
/usr/lib/cloudd/device_info
/usr/lib/cloudd/sync
/usr/lib/cloudd/trafficstatus
/usr/lib/cloudd/runtime
/usr/lib/cloudd/report_topic
/usr/lib/cloudd/time
/usr/lib/cloudd/speed
/usr/lib/cloudd/wifiauth
/usr/lib/lua/cloudd/md5file.lua
/usr/lib/lua/cloudd/firmware.lua
/usr/lib/cloudd/info
/etc/config/cloudd
/usr/lib/cloudd/init
/etc/fixup.d/29-cloudd
/usr/bin/cloudd_cron
/usr/lib/cloudd/reboot
/usr/lib/cloudd/cmd
/usr/lib/cloudd/sms
/etc/uci-defaults/29-cloudd
/usr/lib/cloudd/intstatus
/usr/lib/lua/cloudd/api-new.lua
/usr/lib/lua/cloudd/cloudd.lua
/usr/lib/lua/cloudd/lock.lua
/usr/lib/cloudd/combo
/usr/lib/cloudd/pass
/usr/lib/lua/cloudd/httpfile.lua
