# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/conn_detect.list
/usr/sbin/conn_detect
/etc/config/conn_detect
/etc/init.d/conn_detect
