# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/cpesel.list
/etc/init.d/cpesel
/etc/config/cpesel
/etc/init.d/cpeselgpio
/usr/libexec/rpcd/cpesel1
/usr/libexec/rpcd/cpesel
/usr/bin/cpesel.sh
