# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/factory_siminfo.list
/etc/config/factory_check
/etc/init.d/factory_siminfo
/lib/upgrade/keep.d/factory_siminfo
/usr/sbin/factory_siminfo.sh
