# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/fanctrl.list
/etc/config/fanctrl
/usr/bin/fanctrl.sh
/etc/init.d/fanctrl
