# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/fota_upgrade.list
/usr/bin/fota_upgrade.sh
/etc/init.d/fota_upgrade
/etc/config/fota_upgrade
