# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/fwdd.list
/usr/bin/fwdd
/etc/init.d/fwdd
/etc/config/fwdd
