# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/hostapd-common.list
/etc/wpad/template/hostapd_ra0_wpa3sae.conf
/etc/wpad/template/hostapd_24G.conf
/etc/wpad/template/hostapd_rai0_open_5G.conf
/etc/capabilities/wpad.json
/etc/wpad/template/hostapd_common.conf
/etc/wpad/template/hostapd_rai0_open.conf
/etc/wpad/template/wpa_supplicant_apcli0.conf
/etc/wpad/template/hostapd_rai0_wpa3sae.conf
/usr/bin/hostapd_conf_generate.sh
/etc/wpad/template/hostapd_ra0_wpa2psk.conf
/etc/wpad/template/hostapd_rai0_wpa3sae_5G.conf
/etc/wpad/template/hostapd.logrotate.conf
/usr/share/acl.d/wpad_acl.json
/lib/netifd/hostapd.sh
/etc/init.d/wpad
/etc/rc.button/wps
/lib/netifd/dhcp-get-server.sh
/etc/wpad/template/hostapd_5G.conf
/etc/wpad/template/hostapd_rax0_wpa3sae_6G.conf
/etc/wpad/template/wpa_supplicant_apclix0.conf
/etc/wpad/template/wpa_supplicant.logrotate.conf
/etc/config/wpad
/etc/wpad/template/hostapd_ra0_open.conf
/etc/wpad/template/hostapd_rai0_wpa2psk_5G.conf
/etc/wpad/template/hostapd_rai0_wpa2psk.conf
/etc/wpad/template/hostapd_6G.conf
/etc/wpad/template/wpa_supplicant_apclii0.conf
