# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/htpdate.list
/etc/init.d/htpdate
/usr/sbin/htpdate
/etc/default/htpdate
