# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/infocd.list
/usr/sbin/infocd_consumer
/usr/sbin/infocd
/etc/init.d/infocd
