# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/kpcpe.list
/usr/bin/kpcpe-ctl
/usr/sbin/kpcped
/etc/init.d/kpcped
