# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/kpsh.list
/usr/sbin/kpshd
/etc/uci-defaults/02-kpsh
/etc/init.d/kpsh
