# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/ledctrl.list
/usr/libexec/rpcd/ledctrl
/etc/init.d/ledctrl
/etc/config/ledctrl
/usr/bin/ledctrl.sh
