# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/logservice.list
/usr/sbin/logservice
/usr/sbin/logclient
/etc/init.d/logservice
/etc/config/logservice
