# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/mosquitto-nossl.list
/etc/init.d/mosquitto
/lib/upgrade/keep.d/mosquitto-nossl
/usr/sbin/mosquitto
/etc/config/mosquitto
/usr/local/bin/mosquitto_collect_sysinfo
/etc/mosquitto/mosquitto_acl.conf
