# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/mqttagent.list
/etc/config/mqttagent
/etc/init.d/mqttagent
/usr/bin/mqttagent
