# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/msad.list
/usr/sbin/msad
/etc/uci-defaults/40-msad
/etc/init.d/msad
/etc/fixup.d/30-msad
/etc/config/msad
