# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/nrswitch.list
/usr/libexec/rpcd/nrswitch
/etc/hotplug.d/gmac/02_nrswitch
/etc/init.d/nrswitch
/usr/sbin/nrswitch
