# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/report_proactively.list
/usr/sbin/report_proactively
/etc/init.d/report_proactively
/etc/config/report_proactively
