# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/simcom_http.list
/etc/init.d/simcom_http
/usr/sbin/simcom_http_cli
/usr/sbin/simcom_http
