# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/smsd.list
/usr/sbin/smsd
/etc/init.d/smsd
/etc/config/smsd
