# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/terminal_trackd.list
/etc/config/terminal_trackd
/etc/fixup.d/92-terminal_trackd
/etc/init.d/terminal_trackd
/usr/sbin/terminal_trackd
