# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/wanchk.list
/usr/bin/wanchk.sh
/etc/init.d/wanchk
/etc/fixup.d/90-wanchk
/etc/config/wanchk
/usr/libexec/rpcd/wanchk
