# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/wanswd.list
/etc/init.d/wanswd
/usr/bin/wanswd.sh
/lib/upgrade/keep.d/wanswd
