# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/usr/lib/opkg/info/wm.list
/usr/sbin/wm
/etc/init.d/wm
/etc/config/wm
/usr/lib/lua/wm/update_apcli
