# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/init.d/cellular_power
#!/bin/sh /etc/rc.common

START=05
STOP=75

start() {
	cpetools.sh -k &
}

shutdown() {
	cpetools.sh -d 
}
