# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/init.d/system_power
#!/bin/sh /etc/rc.common

START=05
STOP=75

start() {
	 echo "start"
}

shutdown() {
	if ubus list|grep -q logservice ;then
		ubus call logservice sync
		sync
		sleep 1
	fi
}
