# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/uci-defaults/40_luci-mcast
#!/bin/ash

if ! uci -q get luci.module.mcast >/dev/null; then
	uci set luci.module=internal
	uci set luci.module.mcast=0
	uci commit luci.module
fi

rm -f /tmp/luci-indexcache
exit 0
