# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/uci-defaults/89-ptype
#!/bin/sh

gWdir="/etc/ptype.d"
gType="oem.board.ptype"
gName=

if ! uci -q get "$gType"; then
	gName=$(grep -ioE "^[a-z0-9]+" $gWdir/name)
	uci set "$gType"="${gName:-rt}"
	uci commit oem
fi
