# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/uci-defaults/95-fixup
#!/bin/sh

local _name="fixup"
local _dir="/etc/${_name}.d"

if [ -d "$_dir" ]; then
	find "$_dir" -type f -exec sh {} \;
fi

exit 0
