# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/uci-defaults/99-apcli
#!/bin/sh

uci -q get wireless.apcli &>/dev/null && exit 0

uci -q batch <<-EOF >/dev/null
	set wireless.apcli=wifi-iface
	set wireless.apcli.disabled=1
	set wireless.apcli.device=radio0
	set wireless.apcli.mode=sta
	commit wireless
EOF

exit 0
