# 10-router-pages/factory-bin-C2000MAX-2.1.7.n0.c1-extracted/etc/uci-defaults/luci-theme-nradio
#!/bin/sh
uci batch <<-EOF
	set luci.themes.NRadio=/luci-static/nradio
	set luci.main.mediaurlbase=/luci-static/nradio
	commit luci
EOF
exit 0
